# qwen-api.service
[Unit]
Description=Qwen Image API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/videoapp/qwenApi
EnvironmentFile=/root/videoapp/qwenApi/.env
ExecStart=/root/videoapp/qwenApi/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
